# TradingBot/scripts/systemd/roostoo-bot.service
[Unit]
Description=Roostoo Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/path/to/roostoo-sharpe-bot
Environment="PATH=/path/to/roostoo-sharpe-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/path/to/roostoo-sharpe-bot/venv/bin/python /path/to/roostoo-sharpe-bot/run.py --mode enhanced --once
Restart=on-failure
RestartSec=30
StandardOutput=append:/path/to/roostoo-sharpe-bot/logs/systemd.log
StandardError=append:/path/to/roostoo-sharpe-bot/logs/systemd-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/path/to/roostoo-sharpe-bot/data /path/to/roostoo-sharpe-bot/logs

[Install]
WantedBy=multi-user.target
